# network/if-up.d/static-arp
#! / Bin / sh
arp -i eth0 -s 193.191.177.6 00: XX: 0C: XX: DD: C1
